# service/entrypointClamAV.service
[Unit]
Description=scan clamAV

[Service]
ExecStart=/usr/bin/entrypointClamAV.sh

[Install]
WantedBy=multi-user.target
